# package/debian/filebrowser-quantum.service
[Unit]
Description=FileBrowser Quantum
Documentation=https://filebrowserquantum.com/en/docs/
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/filebrowser-quantum -c /etc/filebrowser.yml
Restart=always
RestartSec=10
User=filebrowser
Group=filebrowser
ProtectHome=true
ProtectSystem=strict
StateDirectory=filebrowser/ filebrowser/sources/
WorkingDirectory=/var/lib/filebrowser

[Install]
WantedBy=multi-user.target
